# board/common/rootfs_overlay/etc/cron.d/ntpsync
# Sync every hour with pool servers
0 * * * * /usr/bin/ntpdate 0.pool.ntp.org 2.pool.ntp.org ntp1.fau.de

